# xiNAS-MCP/xinas-mcp.service
[Unit]
Description=xiNAS MCP Server (Streamable HTTP / SSE / stdio)
Documentation=https://github.com/xinnor/xiNAS
After=network-online.target xinas-nfs-helper.service
Wants=network-online.target xinas-nfs-helper.service

[Service]
Type=simple
ExecStart=/usr/local/bin/xinas-mcp
Restart=on-failure
RestartSec=5s

# Run as root: reads /etc/xinas-mcp/config.json (0640), writes
# /var/log/xinas/mcp-audit.jsonl, talks to /run/xinas-nfs-helper.sock,
# and execs helper subprocesses (xicli, exportfs).
User=root
Group=root

# Stdio transport sees an immediate EOF; HTTP/SSE keep the event loop alive.
StandardInput=null
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xinas-mcp

ProtectHome=true

[Install]
WantedBy=multi-user.target
