# xiNAS-MCP/xinas-api.service
# xinas-api.service — Phase 0 skeleton
#
# Unit for the new xinas-api process introduced by ADR-0001 and
# ADR-0002. Runs alongside xinas-mcp.service (MCP untouched in
# Phase 0; convergence is WS12).
#
# Per ADR-0002 §Hardening, xinas-api runs UNPRIVILEGED — no root,
# no capabilities. The agent (a separate unit, separate PR) is the
# only thing that runs as root. A static system user (xinas-api,
# primary group xinas-admin) owns the SQLite database and audit
# JSONL so the files persist with a stable owner across restarts
# — DynamicUser=yes would have rotated the uid on every start and
# left previously-written state owned by a uid the next instance
# couldn't write to. The rest of the hardening stack (empty
# CapabilityBoundingSet, NoNewPrivileges, ProtectSystem=strict,
# SystemCallFilter@system-service, etc.) is what ADR-0002 actually
# requires and is preserved.
#
# The xinas_api Ansible role creates the xinas-api user + xinas-admin
# group, pre-creates the writable paths (/var/lib/xinas/state,
# /var/log/xinas, /run/xinas) via /usr/lib/tmpfiles.d/xinas-api.conf,
# templates config.json, and copies this unit verbatim. ReadWritePaths
# below grants write access under ProtectSystem=strict.

[Unit]
Description=xiNAS API (Phase 0 — REST + state store)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Re-assert the tmpfiles.d ownership/mode of the writable dirs on EVERY start.
# systemd-tmpfiles-setup runs these only at boot; a plain `systemctl restart`
# does not. If /var/log/xinas or /var/lib/xinas/state ever drift from
# `xinas-api:xinas-admin` (e.g. recreated as root by a partial reinstall, or a
# log-rotation misstep), the unprivileged api can no longer open audit.jsonl /
# the SQLite db and restart-loops with EACCES until someone manually re-runs
# tmpfiles. The `+` prefix runs this as root (before the sandbox), so it can
# chown; it is idempotent and cheap. This makes restart self-healing.
ExecStartPre=+/usr/bin/systemd-tmpfiles --create /usr/lib/tmpfiles.d/xinas-api.conf
ExecStart=/usr/bin/node /opt/xiNAS/xiNAS-MCP/dist/api-server.js
Environment=XINAS_API_CONFIG=/etc/xinas-api/config.json
Restart=on-failure
RestartSec=5s

# Per ADR-0002 §Hardening: no root, no capabilities.
# Static system user — see header comment. The Ansible role creates
# the user with `useradd --system --no-create-home --shell /usr/sbin/nologin`
# and the group with `groupadd --system`.
User=xinas-api
Group=xinas-admin
# Membership in xinas-api grants connect to /run/xinas/agent.sock (0660
# root:xinas-api per ADR-0002). The group has no human members; the
# api user joins it as supplementary so the effective gid stays
# xinas-admin and files created by the api remain operator-readable.
SupplementaryGroups=xinas-api
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Writable storage — the xinas_api Ansible role pre-creates
# /var/lib/xinas/state and /var/log/xinas as xinas-api:xinas-admin
# 0750 via /usr/lib/tmpfiles.d/xinas-api.conf. ReadWritePaths
# grants write access under ProtectSystem=strict.
ReadWritePaths=/var/lib/xinas/state /var/log/xinas

# /run/xinas is created by the xinas_api Ansible role's tmpfiles.d
# entry as root:xinas-admin 0770. systemd-tmpfiles-setup.service
# recreates it on every boot (/run is tmpfs).

[Install]
WantedBy=multi-user.target
