# xiNAS-MCP/xinas-api.service
# xinas-api.service — Phase 0 skeleton
#
# Unit for the new xinas-api process introduced by ADR-0001 and
# ADR-0002. Runs alongside xinas-mcp.service (MCP untouched in
# Phase 0; convergence is WS12).
#
# Per ADR-0002 §Hardening, xinas-api runs UNPRIVILEGED — no root,
# no capabilities. The agent (a separate unit, separate PR) is the
# only thing that runs as root. DynamicUser=yes gives the process a
# transient user; StateDirectory/LogsDirectory give it writable
# storage automatically.
#
# Full Ansible wiring (config.json templating, group setup for the
# UDS, ReadWritePaths tuning) is a follow-up PR — this unit is the
# minimum viable version that starts the service safely.

[Unit]
Description=xiNAS API (Phase 0 — REST + state store)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/xiNAS/xiNAS-MCP/dist/api-server.js
Environment=XINAS_API_CONFIG=/etc/xinas-api/config.json
Restart=on-failure
RestartSec=5s

# Per ADR-0002 §Hardening: no root, no capabilities.
DynamicUser=yes
# The DynamicUser is added to xinas-admin as a supplementary group so
# the api process can chown(2) the UDS to that group at startup. The
# numeric gid is templated into /etc/xinas-api/config.json (listen.socketGroup)
# by the Ansible role (deferred follow-up — see PR #201 'What's deferred').
# Without that config value, the socket is owned by the dynamic user's
# primary group and only the api process itself can connect — safe but
# unreachable from operator tools.
SupplementaryGroups=xinas-admin
NoNewPrivileges=true
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Writable storage — the xinas_api Ansible role pre-creates
# /var/lib/xinas/state and /var/log/xinas as root:xinas-admin
# 0770 via /usr/lib/tmpfiles.d/xinas-api.conf. ReadWritePaths
# grants the DynamicUser write access under ProtectSystem=strict.
ReadWritePaths=/var/lib/xinas/state /var/log/xinas

# /run/xinas is created by the xinas_api Ansible role's tmpfiles.d
# entry as root:xinas-admin 0750 (NOT via RuntimeDirectory, which
# under DynamicUser would assign ownership to the transient uid and
# block xinas-admin traversal). systemd-tmpfiles-setup.service
# recreates it on every boot.

[Install]
WantedBy=multi-user.target
