# xiNAS-MCP/nfs-helper/xinas-nfs-helper.service
[Unit]
Description=xiNAS NFS Helper Daemon
Documentation=https://github.com/xinnor/xiNAS
After=network.target nfs-kernel-server.service
Wants=nfs-kernel-server.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/xinas-mcp/nfs-helper/nfs_helper.py
Restart=on-failure
RestartSec=5s

# Runtime directory creates /run/xinas-nfs-helper/
RuntimeDirectory=xinas-nfs-helper
RuntimeDirectoryMode=0755

# Run as root (required for /etc/exports and xfs_quota)
User=root
Group=root

# Minimal hardening — service runs as root and manages /etc/exports, /run
# ProtectSystem=full is intentionally omitted: it makes /etc read-only and
# ReadWritePaths can only re-open paths that already exist at service start,
# so a missing /etc/exports would remain read-only (PermissionError on create).
ProtectHome=true

Environment=NFS_HELPER_SOCKET=/run/xinas-nfs-helper.sock
Environment=LOG_LEVEL=INFO

StandardOutput=journal
StandardError=journal
SyslogIdentifier=xinas-nfs-helper

[Install]
WantedBy=multi-user.target
